# roles/daily-borgbackup/files/cron.daily/local-borgbackup
#!/bin/sh
# Tested on Stretch i.e. borgbackup-1.0.9

BORG_REPOS_DIR=/data/backup/borg

# Show details when run from a terminal (idea stolen from bup).
# Otherwise, only show errors.
BORG_PROGRESS=
BORG_STATS=
if tty>/dev/null; then
        BORG_PROGRESS="--progress"
        BORG_STATS="--stats --info"
fi

# Write a log to show the time taken
log() {
        echo >> /var/log/local-borgbackup -- "$@"
}
log_date() {
        log "$(date --rfc-3339=seconds)"
}
log_cmd() {
        "$@"
        log "$*"
        log_date  # log time of completion
}
# Log the start timestamp
log_date

# If the borg repo does not exist, create it - atomically.
# If interrupted, just leaves a small randomly-named directory lying around
borg_init() {
	if [ ! -e "$1" ]; then
		T=$(mktemp --dry-run "$1".tmp.XXXXXX)
		borg init --encryption=none "$T"

		# ~5M segment size (borg 1.0 default) is too small IMO.
		# Use ~500M (borg 1.1 default).
		sed -i "$T"/config -e 's/^max_segment_size = .*$/max_segment_size = 524288000/'

		mv "$T" "$1"
	fi

	# borg doesn't like renames.
	# There's an easier fix in 1.1, which might even work before the rename
	# (so we don't end up overriding borg *after* the initial creation).
	# It would also avoid a warning which we have to suppress.
	# https://github.com/borgbackup/borg/issues/1910#issuecomment-354780246
	(
		export BORG_REPO="$1"
		export BORG_RELOCATED_REPO_ACCESS_IS_OK=yes
		borg delete ::temp 2>/dev/null
		borg create ::temp - </dev/null 2>/dev/null
		borg delete ::temp
	)
}

backup_home() {
	HOME_DIR=/home/"$1"
	BORG_REPO="$BORG_REPOS_DIR"/"$1"
	export BORG_REPO

	borg_init "$BORG_REPO"

	# Exclusions:
	#
	# Firefox profile tends to churn a whole lot (but see below).
	#
	# Thunderbird has an index for searching messages, which shows up.
	# This is due to the size of the blocks borg uses (megabytes?),
	# although enabling compression helps a bit.
	# We're allowed to drop the index - it will be rebuilt.
	#
	# On Debian, there may be an icedove profile which is half-migrated
	# to thunderbird, but still uses the old location.
	#
	# Nowadays GNOME logs activity and indexes it, that shows up as well.
	#
	log_cmd borg create $BORG_PROGRESS $BORG_STATS '::{now}' -C lz4 \
		-x "$HOME_DIR" \
		--exclude "$HOME_DIR"/.cache \
		--exclude "$HOME_DIR"/.mozilla/firefox \
		--exclude "$HOME_DIR"/.thunderbird/*/global-messages-db.sqlite \
		--exclude "$HOME_DIR"/.icedove/*/global-messages-db.sqlite \
		--exclude "$HOME_DIR"/.local/share/zeitgeist/fts.index \
		--exclude "$HOME_DIR"/.local/share/tracker

	# Prune old backups to save space
	log_cmd borg prune $BORG_STATS \
		--keep-daily 14 --keep-weekly 4 \
		--keep-monthly 12 --keep-yearly -1


	# Also save bookmarkbackups from the firefox profile.
	# Firefox Sync doesn't e.g. protect against malice,
	# and I have a lot of bookmarks I want to keep.
	if [ -d "$HOME_DIR"/.mozilla ]; then
		BORG_REPO="$BORG_REPOS_DIR"/"$1".firefox-bookmarks
		export BORG_REPO

		borg_init "$BORG_REPO"

		log_cmd borg create $BORG_PROGRESS $BORG_STATS '::{now}' "$HOME_DIR"/.mozilla/firefox/*/bookmarkbackups

		log_cmd borg prune $BORG_STATS \
			--keep-daily 14 --keep-weekly 4 \
			--keep-monthly 12 --keep-yearly -1
	fi
}

backup_home mike
backup_home sheena
backup_home jenkins-photos
