# roles/daily-borgbackup/files/cron.daily/local-borgbackup
#!/bin/sh
# Tested on Stretch i.e. borgbackup-1.0.9
#
# TODO: as part of init
# To upgrade configs to borg 1.1 segment sizes:
# 
# sed -i "$BORG_REPO"/config -e 's/^max_segment_size = 5242880$/max_segment_size = 524288000/'

# Show progress when run from a terminal (idea stolen from bup)
BORG_PROGRESS=
BORG_STATS=
if tty>/dev/null; then
        BORG_PROGRESS="--progress"
        BORG_STATS="--stats --info"
fi


log() {
        echo >> /var/log/local-borgbackup -- "$@"
}
log_date() {
        log "$(date --rfc-3339=seconds)"
}

# Log start timestamp
log_date

log_cmd() {
        "$@"
        log "$*"
        log_date  # log time of completion
}

cd /data/backup/borg || exit 1
for i in *; do
	[ -d "$i" ] || continue

	HOME_DIR=/home/"$i"
	export BORG_REPO="$PWD"/"$i"

	# Exclusions.
	#
	# Firefox profile tends to churn a whole lot (but see below).
	#
	# Thunderbird has an index for searching messages, which shows up.
	# This is due to the size of the blocks borg uses (megabytes?),
	# although enabling compression helps a bit.
	# We're allowed to drop the index - it will be rebuilt.
	#
	# On Debian, there may be an icedove profile which is half-migrated
	# to thunderbird, but still uses the old location.
	#
	# Nowadays GNOME logs activity and indexes it, that shows up as well.
	#
	log_cmd borg create $BORG_PROGRESS $BORG_STATS '::home {now}' -C lz4 \
		-x "$HOME_DIR" \
		--exclude "$HOME_DIR"/.cache \
		--exclude "$HOME_DIR"/.mozilla/firefox \
		--exclude "$HOME_DIR"/.thunderbird/*/global-messages-db.sqlite \
		--exclude "$HOME_DIR"/.icedove/*/global-messages-db.sqlite \
		--exclude "$HOME_DIR"/.local/share/zeitgeist/fts.index \
		--exclude "$HOME_DIR"/.local/share/tracker

	# Add bookmarkbackups from the firefox profile.
	# Firefox Sync doesn't e.g. protect against malice,
	# and I have a lot of bookmarks I want to keep.
	if [ -d "$HOME_DIR"/.mozilla ]; then
		log_cmd borg create $BORG_PROGRESS $BORG_STATS '::firefox-bookmarks {now}' "$HOME_DIR"/.mozilla/firefox/*/bookmarkbackups
	fi
done

for i in *; do
	[ -d "$i" ] || continue
	export BORG_REPO="$PWD"/"$i"

	log_cmd borg prune $BORG_STATS --prefix "home " \
		--keep-daily 14 --keep-weekly 4 \
		--keep-monthly 12 --keep-yearly -1

	log_cmd borg prune $BORG_STATS --prefix "firefox-bookmarks " \
		--keep-daily 14 --keep-weekly 4 \
		--keep-monthly 12 --keep-yearly -1
done
